# scripts/deb/DEBIAN/postinst
#!/bin/bash

if [ -e /tmp/hautomation-config.json ]
then
    rm -Rf /var/hautomation/data/config.json
    mv /tmp/hautomation-config.json /var/hautomation/data/config.json
fi

systemctl daemon-reload
systemctl enable hautomation
systemctl restart hautomation
echo ""
echo "\033[32mHautomation has been successfully installed !\033[0m"
echo ""
echo "\033[4mHow to use ?\033[0m"
echo "Open a web browser on your computer and go to \033[36mhttps://me.hautomation-io.com/\033[0m to start, with the username \033[33madmin\033[0m and password \033[33madmin\033[0m ! (on first installation, start can take some times to install dependencies, wait for some minutes if it dosen't works)"
echo ""
